# supervisor/default.conf
[supervisord]
user=api
nodaemon=true
stdout_logfile=/dev/stdout
stdout_logfile_maxbytes = 0
stderr_logfile=/dev/stderr
stderr_logfile_maxbytes=0

[program:backend]
user=api
process_name=%(program_name)s_%(process_num)02d
command=uwsgi --http :5000 --wsgi-file /app/wsgi.py --workers 4 --threads 4 -w wsgi:app

directory=/app/
stdout_logfile=/dev/stdout
stdout_logfile_maxbytes = 0
stderr_logfile=/dev/stderr
stderr_logfile_maxbytes=0